# br2_external/package/pi4_feature_check/pi4-feature-check.service
[Unit]
Description=Pi4 Feature Checklist Boot Report
After=network-online.target av-core.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/pi4_feature_check.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
