# br2_external/package/av_services/av-core-gateway.service
[Unit]
Description=AV Core Gateway Service
Documentation=file:///etc/av-core/gateway.conf
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/av-core-gateway
EnvironmentFile=-/etc/av-core/gateway.conf
Restart=on-failure
RestartSec=5s
WatchdogSec=30s
NotifyAccess=main
StandardOutput=journal
StandardError=journal
SyslogIdentifier=av-core-gateway
StartLimitIntervalSec=60s
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
